# llm-api.service
[Unit]
Description=Internal LLM API Server
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/llm-api
Environment="PATH=/opt/llm-api/venv/bin"
ExecStart=/opt/llm-api/venv/bin/uvicorn src.main:app --host 0.0.0.0 --port 8000 --workers 4
Restart=always
RestartSec=10
StandardOutput=append:/var/log/llm-api/stdout.log
StandardError=append:/var/log/llm-api/stderr.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/llm-api/logs

[Install]
WantedBy=multi-user.target
